# desktop/.config/systemd/user/tailscale-systray.service
[Unit]
Description=Tailscale System Tray
After=graphical.target

[Service]
Type=simple
ExecStart=tailscale systray

[Install]
WantedBy=default.target
