# rootfs/etc/udev/rules.d/99-removable.rules
ACTION=="add", SUBSYSTEM=="thunderbolt", ATTR{authorized}=="0", ATTR{authorized}="1", RUN+="/etc/udev/scripts/pcie-rescan.sh"
