# hyprland/.config/systemd/user/swayosd.service
[Unit]
Description=SwayOSD server for on-screen display notifications
Documentation=https://github.com/ErikReider/SwayOSD
PartOf=hyprland-session.target
After=hyprland-session.target

[Service]
Type=simple
ExecStart=/usr/bin/swayosd-server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=hyprland-session.target
