# waybar/.config/systemd/user/waybar@sway.service.d/override.conf
[Unit]
ConditionEnvironment=XDG_CURRENT_DESKTOP=sway
PartOf=
PartOf=graphical-session.target
After=
After=graphical-session.target

[Install]
WantedBy=
WantedBy=graphical-session.target
