# waybar/.config/systemd/user/waybar.service
[Unit]
Description=Waybar - Highly customizable Wayland bar (generic/unused)
Documentation=https://github.com/Alexays/Waybar
PartOf=graphical-session.target
After=graphical-session.target
Before=xdg-desktop-autostart.target
Requires=dbus.socket

[Service]
Type=simple
ExecCondition=/bin/sh -c 'command -v waybar'
ExecStart=waybar
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
