# waybar/.config/systemd/user/waybar.service
# UNUSED - This service file is not used
# Use waybar@.service template instead with compositor-specific instances:
#   - waybar@sway.service for Sway
#   - waybar@hyprland.service for Hyprland
#
# This file is kept as an example of generic waybar service configuration
[Unit]
Description=Waybar - Highly customizable Wayland bar (generic/unused)
Documentation=https://github.com/Alexays/Waybar
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecCondition=/bin/sh -c 'command -v waybar'
ExecStart=waybar
Restart=on-failure
RestartSec=3
# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
