# waybar/.config/systemd/user/waybar@.service
[Unit]
Description=Waybar - Highly customizable Wayland bar for %i
Documentation=https://github.com/Alexays/Waybar
PartOf=%i-session.target
After=%i-session.target

[Service]
Type=simple
ExecCondition=/bin/sh -c 'command -v waybar'
ExecStart=waybar -c %h/.config/waybar/%i.jsonc
# Hold "active" until waybar publishes org.kde.StatusNotifierWatcher,
# so units that After=waybar@%i.service (e.g. xdg-desktop-autostart.target
# via its drop-in) actually wait for the tray to be ready instead of
# just for the waybar process to exec.
ExecStartPost=/bin/sh -c 'for i in $(seq 1 50); do busctl --user --no-pager list 2>/dev/null | grep -q org.kde.StatusNotifierWatcher && break; sleep 0.1; done'
Restart=on-failure
RestartSec=3
# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=%i-session.target
