# hyprland/.config/systemd/user/hyprpaper.service.d/override.conf
[Service]
UnsetEnvironment=AQ_DRM_DEVICES MANGOHUD
Environment=DISABLE_MANGOHUD=1
