# hyprland/.config/systemd/user/hyprpaper.service.d/override.conf
[Service]
UnsetEnvironment=MANGOHUD
Environment=DISABLE_MANGOHUD=1
