# hyprwhspr/.config/systemd/user/hyprwhspr.service.d/override.conf
[Service]
ExecStart=
ExecStart=/bin/bash -c '[ -r %h/.zshsecrets ] && . %h/.zshsecrets; exec /usr/lib/hyprwhspr/bin/hyprwhspr'
