# rootfs/etc/systemd/system/clamav-clamonacc.service.d/override.conf
[Service]
ExecStart=clamonacc -F --fdpass --log=/var/log/clamav/clamonacc.log
