# xsettingsd/.config/systemd/user/xsettingsd.service
[Unit]
Description=XSettings daemon for Hyprland XWayland clients
Documentation=man:xsettingsd(1)
PartOf=graphical-session.target
After=wayland-session-waitenv.service hyprland-theme.service
Before=graphical-session.target xdg-desktop-autostart.target
ConditionEnvironment=DISPLAY

[Service]
Type=simple
ExecCondition=/bin/sh -c 'command -v xsettingsd'
ExecStart=xsettingsd -c %h/.config/xsettingsd/xsettingsd.conf
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
