# hyprpilot/.config/systemd/user/hyprpilot.service
[Unit]
Description=hyprpilot ACP overlay daemon (temporary)
Documentation=https://gitlab.kilic.dev/utils/hyprpilot
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecCondition=/bin/sh -c 'command -v hyprpilot'
Environment=RUST_LOG=warn,hyprpilot=info,acp::emit=trace,tauri::emit=trace,snapshot::mirror=trace,snapshot::meta=trace,snapshot::chat=trace,acp::wire=trace,acp::thought=trace,webview=trace
ExecStart=/bin/zsh -c '%h/.local/bin/hyprpilot daemon --hidden'
Restart=on-failure
RestartSec=3
# Surface daemon stdout / stderr in `journalctl --user -u hyprpilot`.
# The daemon also writes its own rolling file log under
# `$XDG_STATE_HOME/hyprpilot/logs/` — both paths are useful while
# this binary is shipped as a dropped-in /usr/local/bin/.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
